# vts-exporter.ini
[program:vts-exporter]
command=/bin/nginx-vts-exporter -nginx.scrape_uri=http://localhost:8080/status/format/json
stdout_events_enabled=true
stderr_events_enabled=true
autostart = true
autorestart = true

